# etc/cron.d/avidbots-prometheus-adapter
*/10 * * * * root /usr/bin/python /etc/avidbots-prometheus-adapter/prometheus_adapter.py
